# deploy/bingxi.service
[Unit]
Description=Bingxi ERP Backend Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/bingxi/backend
EnvironmentFile=/etc/bingxi/.env
ExecStart=/opt/bingxi/backend/server
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bingxi

# 安全限制
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/bingxi/logs /var/log/bingxi

[Install]
WantedBy=multi-user.target
