# deploy/bingxi-backend.service
[Unit]
Description=Bingxi ERP Backend Service
After=network.target postgresql.service redis.service
Wants=network-online.target

[Service]
Type=simple
User=bingxi
Group=bingxi
WorkingDirectory=/opt/bingxi-erp/backend
EnvironmentFile=/etc/bingxi-erp/.env
Environment=CONFIG_FILE=/opt/bingxi-erp/backend/config.yaml
ExecStart=/opt/bingxi-erp/backend/server
Restart=always
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60

LimitNOFILE=65535

StandardOutput=journal
StandardError=journal
SyslogIdentifier=bingxi-backend

[Install]
WantedBy=multi-user.target
